# dotfiles/dotfiles/etc/systemd/system/mnt-windows.automount
# Windows shared filesystem systemd automount unit

[Unit]
Description=Windows shared filesystem automount

[Install]
WantedBy=multi-user.target

[Automount]
Where=/mnt/windows
