# dotfiles/dotfiles/etc/systemd/system/mnt-vm-shared.automount
# Virtual machine shared directories systemd automount unit

[Unit]
Description=Virtual machine shared directories automount

[Install]
WantedBy=multi-user.target

[Automount]
Where=/mnt/vm/shared
