# dotfiles/dotfiles/etc/systemd/system/mnt-raspberry.automount
# Raspberry filesystem systemd automount unit

[Unit]
Description=Raspberry filesystem automount

[Install]
WantedBy=multi-user.target

[Automount]
Where=/mnt/raspberry
