# dotfiles/config/systemd/user/shell-cache.timer
# shell-cache systemd timer file

[Unit]
Description=Job refreshing shell caches

[Install]
WantedBy=multi-user.target

[Timer]
OnCalendar=*-*-* 00:05:00
Persistent=true
